# haproxy-reload/sv/reloader/run
#!/bin/bash

# reload when config changes
while inotifywait -e close_write /etc/haproxy/haproxy.cfg; do
	echo "haproxy config changed"
  	sv reload haproxy
done